# debian/DEBIAN/prerm
#!/bin/bash

PKGNAME="__PKGNAME__"

# Останавливаем и отключаем сервис
systemctl stop $PKGNAME.service
systemctl disable $PKGNAME.service

exit 0 