# debian/DEBIAN/prerm
#!/bin/bash

# Останавливаем и отключаем сервис
systemctl stop todoist-bot.service
systemctl disable todoist-bot.service

exit 0 